# etc/freebsd/okld.sh.in
#!/bin/sh
#
# OKWS rc.d control script
#
# This is a standard rc.d script to control okws.  It responds to all
# standard rc.d commands:
#
#   Usage: okld.sh [fast|force](start|stop|restart|rcvar|status|poll)
#
# In order for it to work, okld_enable="YES" needs to be set in 
# /etc/rc.conf.
#
# I don't recommend using 'fast'.
#
# Author:  Patrick Crosby <patrick@okcupid.com>
#
# $Id$
#

. /etc/rc.subr

name="okld"
rcvar=`set_rcvar`
command="%%PREFIX%%/%%EXE_BIN%%/okld"
okld_flags="-q"
pidfile="/var/run/okld.pid"
required_files="%%PREFIX%%/etc/okws/okws_config"
extra_command="reload"

load_rc_config $name
run_rc_command "$1"

